# linux-deployment/supervisord.conf
[program:nvidia-proxy]
command=/usr/bin/node /opt/nvidia-to-openai-proxy/server.js
directory=/opt/nvidia-to-openai-proxy
user=nvidia-proxy
autostart=true
autorestart=true
startsecs=10
startretries=3
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nvidia-proxy-out.log
stderr_logfile=/var/log/supervisor/nvidia-proxy-err.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=NODE_ENV="production",PORT="3000",NVIDIA_API_KEY="your-api-key-here"
